# roles/base/templates/cupsresume.service.j2
[Unit]
Description=Restart cups after resume
After=suspend.target
[Service]
Type=simple
ExecStart=/usr/bin/systemctl restart cups
[Install]
WantedBy=suspend.target